# reposerver_service_script
#!/bin/sh
### BEGIN INIT INFO
# Provides:          reposerver
# Required-Start:    $remote_fs $network
# Required-Stop:     $remote_fs $network
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start and stop the Reposerver Platform
# Description:       High-performance hybrid Python/C++ application server.
### END INIT INFO

# ====================[ CONFIGURAÇÃO ]====================
APP_DIR="/opt/reposerver"
APP_USER="reposerver"
NAME="reposerver"
VENV_DIR="$APP_DIR/venv"
GUNICORN_EXEC="$VENV_DIR/bin/gunicorn"
APP_MODULE="reposerver_main:app"
PIDFILE="/var/run/$NAME.pid"
LOG_DIR="$APP_DIR/logs"
GUNICORN_LOGFILE="$LOG_DIR/gunicorn.log"
DEBUG_LOGFILE="/tmp/reposerver_debug.log"
WORKERS=3
WORKER_CLASS="eventlet"
HOST="0.0.0.0"
PORT="5000"
GUNICORN_ARGS="--workers $WORKERS --worker-class $WORKER_CLASS --bind $HOST:$PORT --log-file $GUNICORN_LOGFILE"

. /lib/lsb/init-functions

# ====================[ FUNÇÕES DO SERVIÇO ]====================

prepare_dirs() {
    echo "DEBUG: Running prepare_dirs" >> $DEBUG_LOGFILE
    mkdir -p "$LOG_DIR"
    echo "DEBUG: Executing chown on $LOG_DIR" >> $DEBUG_LOGFILE
    chown -R "$APP_USER":"$APP_USER" "$LOG_DIR"
    ls -ld "$LOG_DIR" >> $DEBUG_LOGFILE
    
    mkdir -p /var/run
    echo "DEBUG: Executing chown on /var/run" >> $DEBUG_LOGFILE
    chown -R "$APP_USER":"$APP_USER" /var/run
    ls -ld /var/run >> $DEBUG_LOGFILE
}

start() {
    echo "====================[ NEW START ATTEMPT ]====================" > $DEBUG_LOGFILE
    echo "Timestamp: $(date)" >> $DEBUG_LOGFILE
    echo "Running as user: $(whoami)" >> $DEBUG_LOGFILE
    
    echo "
[DEBUG] Checking permissions before starting:" >> $DEBUG_LOGFILE
    ls -ld "$APP_DIR" >> $DEBUG_LOGFILE
    ls -lR "$APP_DIR" >> $DEBUG_LOGFILE
    
    prepare_dirs

    echo "
[DEBUG] Final check of variables:" >> $DEBUG_LOGFILE
    echo "GUNICORN_EXEC: $GUNICORN_EXEC" >> $DEBUG_LOGFILE
    ls -l "$GUNICORN_EXEC" >> $DEBUG_LOGFILE
    echo "APP_USER: $APP_USER" >> $DEBUG_LOGFILE
    echo "PIDFILE: $PIDFILE" >> $DEBUG_LOGFILE
    echo "GUNICORN_LOGFILE: $GUNICORN_LOGFILE" >> $DEBUG_LOGFILE
    
    if [ -f $PIDFILE ]; then
        log_daemon_msg "$NAME já está em execução (PID: $(cat $PIDFILE))"
        echo "DEBUG: Service already running." >> $DEBUG_LOGFILE
        return 1
    fi
    
    log_daemon_msg "Iniciando $NAME"
    
    cd "$APP_DIR" || { echo "DEBUG: Failed to cd to $APP_DIR" >> $DEBUG_LOGFILE; return 1; }
    
    echo "
[DEBUG] About to execute start-stop-daemon..." >> $DEBUG_LOGFILE
    
    start-stop-daemon --start --quiet --pidfile "$PIDFILE" --make-pidfile \
        --background --chuid "$APP_USER" --chdir "$APP_DIR" --startas /bin/bash -- \
        -c "exec $GUNICORN_EXEC $APP_MODULE $GUNICORN_ARGS >> $GUNICORN_LOGFILE 2>&1"

    sleep 2
    if ! status > /dev/null; then
        log_end_msg 1
        log_failure_msg "Falha ao iniciar. Verifique os logs em $GUNICORN_LOGFILE e o log de depuração em $DEBUG_LOGFILE"
        echo "DEBUG: start-stop-daemon failed." >> $DEBUG_LOGFILE
    else
        log_end_msg 0
        echo "DEBUG: Service started successfully." >> $DEBUG_LOGFILE
    fi
}

stop() {
    if [ ! -f "$PIDFILE" ]; then
        log_daemon_msg "$NAME não está em execução"
        return 1
    fi
    
    log_daemon_msg "Parando $NAME"
    start-stop-daemon --stop --quiet --pidfile "$PIDFILE" --retry 5
    
    if [ -f "$PIDFILE" ]; then
        rm -f "$PIDFILE"
    fi
    
    log_end_msg $?
}

status() {
    status_of_proc -p "$PIDFILE" "$GUNICORN_EXEC" "$NAME"
    return $?
}

case "$1" in
    start)
        start
        ;;
    stop)
        stop
        ;;
    restart)
        log_daemon_msg "Reiniciando $NAME"
        stop
        sleep 1 
        start
        ;;
    status)
        status
        ;;
    *)
        echo "Uso: /etc/init.d/$NAME {start|stop|restart|status}"
        exit 1
        ;;
esac

exit 0
