# deploy/updater.service
[Unit]
Description=Reposerver auto-updater
After=network.target

[Service]
User=%i
WorkingDirectory=%h/Reposerver
ExecStart=/usr/bin/env python3 %h/Reposerver/updater.py
Restart=always
RestartSec=60

[Install]
WantedBy=multi-user.target
